# root/etc/service/fpm/run
#!/bin/sh
exec /usr/sbin/php-fpm7.0 -F -R >>/log/fpm.log 2>&1
